# ops/systemd/openvision-jetson.service
[Unit]
Description=OpenVision Rokid v2 Jetson Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openvision
WorkingDirectory=/opt/openvision-rokid/jetson/agent
Environment=OPENVISION_ENV=jetson
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=OPENVISION_REALTIME_MODEL=gpt-realtime-1.5
Environment=OPENVISION_DEBUG_STT_ENABLED=0
Environment=OPENVISION_DEBUG_STT_LANGUAGE=vi
Environment=OPENVISION_DEBUG_STT_PROFILE=v2_debug_sentence
Environment=OPENVISION_DEBUG_STT_BEAM_SIZE=5
Environment=OPENVISION_DEBUG_STT_VAD_FILTER=1
Environment=OPENVISION_DEBUG_STT_MIN_AUDIO_MS=350
Environment=OPENVISION_DEBUG_STT_MAX_AUDIO_MS=12000
Environment="OPENVISION_DEBUG_STT_HOTWORDS=áo vàng, áo xanh, đếm người, tìm người, nghe được không"
Environment=OPENAI_API_KEY_FILE=/opt/openvision-rokid/ops/secrets/openai_api_key
Environment=OPENVISION_RV101_TCP_INGEST=1
Environment=OPENVISION_RV101_BIND_HOST=0.0.0.0
Environment=OPENVISION_RV101_ADVERTISED_HOST=127.0.0.1
Environment=OPENVISION_RV101_VIDEO_PORT=8770
Environment=OPENVISION_RV101_AUDIO_PORT=8771
EnvironmentFile=-/opt/openvision-rokid/ops/openvision.env
ExecStart=/opt/openvision-rokid/.venv/bin/python -m uvicorn openvision_jetson.fastapi_app:app --host 0.0.0.0 --port 8765 --timeout-graceful-shutdown 3
Restart=on-failure
RestartSec=2
TimeoutStopSec=8

[Install]
WantedBy=multi-user.target
